# ha_opencode/rootfs/etc/s6-overlay/s6-rc.d/ha-opencode/run
#!/command/with-contenv bashio
# ==============================================================================
# HA OpenCode - Main Service
# ==============================================================================

bashio::log.info "=============================================="
bashio::log.info "  HA OpenCode for Home Assistant"
bashio::log.info "  Starting services..."
bashio::log.info "=============================================="

# Set home to persistent data directory
export HOME="/data"
export XDG_DATA_HOME="/data/.local/share"
export XDG_CONFIG_HOME="/data/.config"

# Ensure data directories exist
mkdir -p /data/.local/share/opencode
mkdir -p /data/.config/opencode

# Sync config from template (MCP/LSP enabled state comes from bashio config, not preserved)
if [ ! -f "/data/.config/opencode/opencode.json" ]; then
    cp /opt/ha-mcp-server/opencode-ha.json /data/.config/opencode/opencode.json
    bashio::log.info "Created default OpenCode configuration"
else
    # Update from template but don't preserve enabled states - those come from bashio config
    cp /opt/ha-mcp-server/opencode-ha.json /data/.config/opencode/opencode.json
    bashio::log.info "Updated OpenCode configuration from template"
fi

# ==============================================================================
# Read configuration options
# ==============================================================================

# Feature toggles
MCP_ENABLED=$(bashio::config 'mcp_enabled')

# Terminal appearance
TERMINAL_THEME=$(bashio::config 'terminal_theme')
FONT_SIZE=$(bashio::config 'font_size')
CURSOR_STYLE=$(bashio::config 'cursor_style')
CURSOR_BLINK=$(bashio::config 'cursor_blink')

# LSP feature toggle (always enabled by default - it's read-only)
LSP_ENABLED=$(bashio::config 'lsp_enabled' || echo "true")

bashio::log.info "Configuration:"
bashio::log.info "  - MCP Integration: ${MCP_ENABLED}"
bashio::log.info "  - LSP Integration: ${LSP_ENABLED}"
bashio::log.info "  - Terminal Theme: ${TERMINAL_THEME}"
bashio::log.info "  - Font Size: ${FONT_SIZE}px"

# ==============================================================================
# Apply MCP configuration
# ==============================================================================

if [ "${MCP_ENABLED}" = "true" ]; then
    bashio::log.info "Enabling Home Assistant MCP integration..."
    if command -v jq &> /dev/null; then
        jq '.mcp.homeassistant.enabled = true' "/data/.config/opencode/opencode.json" > "/data/.config/opencode/opencode.json.tmp" && \
            mv "/data/.config/opencode/opencode.json.tmp" "/data/.config/opencode/opencode.json"
    fi
else
    bashio::log.info "MCP integration disabled"
    if command -v jq &> /dev/null; then
        jq '.mcp.homeassistant.enabled = false' "/data/.config/opencode/opencode.json" > "/data/.config/opencode/opencode.json.tmp" && \
            mv "/data/.config/opencode/opencode.json.tmp" "/data/.config/opencode/opencode.json"
    fi
fi

# ==============================================================================
# Apply LSP configuration
# ==============================================================================

if [ "${LSP_ENABLED}" = "true" ]; then
    bashio::log.info "Enabling Home Assistant LSP integration..."
    if command -v jq &> /dev/null; then
        # Remove disabled flag if it exists (LSP is enabled by default if not disabled)
        jq 'if .lsp["ha-yaml"].disabled then del(.lsp["ha-yaml"].disabled) else . end' "/data/.config/opencode/opencode.json" > "/data/.config/opencode/opencode.json.tmp" && \
            mv "/data/.config/opencode/opencode.json.tmp" "/data/.config/opencode/opencode.json"
    fi
else
    bashio::log.info "LSP integration disabled"
    if command -v jq &> /dev/null; then
        # Set disabled: true to disable the LSP server
        jq '.lsp["ha-yaml"].disabled = true' "/data/.config/opencode/opencode.json" > "/data/.config/opencode/opencode.json.tmp" && \
            mv "/data/.config/opencode/opencode.json.tmp" "/data/.config/opencode/opencode.json"
    fi
fi

# ==============================================================================
# Theme definitions
# ==============================================================================

get_theme() {
    local theme_name="$1"
    
    case "${theme_name}" in
        "breeze")
            # KDE Konsole Breeze theme
            echo '{"background":"#232627","foreground":"#fcfcfc","cursor":"#fcfcfc","cursorAccent":"#232627","selectionBackground":"#3daee9","selectionForeground":"#fcfcfc","black":"#232627","red":"#ed1515","green":"#11d116","yellow":"#f67400","blue":"#1d99f3","magenta":"#9b59b6","cyan":"#1abc9c","white":"#fcfcfc","brightBlack":"#7f8c8d","brightRed":"#c0392b","brightGreen":"#1cdc9a","brightYellow":"#fdbc4b","brightBlue":"#3daee9","brightMagenta":"#8e44ad","brightCyan":"#16a085","brightWhite":"#ffffff"}'
            ;;
        "catppuccin_mocha")
            # Catppuccin Mocha theme
            echo '{"background":"#1e1e2e","foreground":"#cdd6f4","cursor":"#f5e0dc","cursorAccent":"#1e1e2e","selectionBackground":"#585b70","selectionForeground":"#cdd6f4","black":"#45475a","red":"#f38ba8","green":"#a6e3a1","yellow":"#f9e2af","blue":"#89b4fa","magenta":"#f5c2e7","cyan":"#94e2d5","white":"#bac2de","brightBlack":"#585b70","brightRed":"#f38ba8","brightGreen":"#a6e3a1","brightYellow":"#f9e2af","brightBlue":"#89b4fa","brightMagenta":"#f5c2e7","brightCyan":"#94e2d5","brightWhite":"#a6adc8"}'
            ;;
        "catppuccin_latte")
            # Catppuccin Latte theme (light)
            echo '{"background":"#eff1f5","foreground":"#4c4f69","cursor":"#dc8a78","cursorAccent":"#eff1f5","selectionBackground":"#acb0be","selectionForeground":"#4c4f69","black":"#5c5f77","red":"#d20f39","green":"#40a02b","yellow":"#df8e1d","blue":"#1e66f5","magenta":"#ea76cb","cyan":"#179299","white":"#acb0be","brightBlack":"#6c6f85","brightRed":"#d20f39","brightGreen":"#40a02b","brightYellow":"#df8e1d","brightBlue":"#1e66f5","brightMagenta":"#ea76cb","brightCyan":"#179299","brightWhite":"#bcc0cc"}'
            ;;
        "dracula")
            # Dracula theme
            echo '{"background":"#282a36","foreground":"#f8f8f2","cursor":"#f8f8f2","cursorAccent":"#282a36","selectionBackground":"#44475a","selectionForeground":"#f8f8f2","black":"#21222c","red":"#ff5555","green":"#50fa7b","yellow":"#f1fa8c","blue":"#bd93f9","magenta":"#ff79c6","cyan":"#8be9fd","white":"#f8f8f2","brightBlack":"#6272a4","brightRed":"#ff6e6e","brightGreen":"#69ff94","brightYellow":"#ffffa5","brightBlue":"#d6acff","brightMagenta":"#ff92df","brightCyan":"#a4ffff","brightWhite":"#ffffff"}'
            ;;
        "nord")
            # Nord theme
            echo '{"background":"#2e3440","foreground":"#d8dee9","cursor":"#d8dee9","cursorAccent":"#2e3440","selectionBackground":"#434c5e","selectionForeground":"#d8dee9","black":"#3b4252","red":"#bf616a","green":"#a3be8c","yellow":"#ebcb8b","blue":"#81a1c1","magenta":"#b48ead","cyan":"#88c0d0","white":"#e5e9f0","brightBlack":"#4c566a","brightRed":"#bf616a","brightGreen":"#a3be8c","brightYellow":"#ebcb8b","brightBlue":"#81a1c1","brightMagenta":"#b48ead","brightCyan":"#8fbcbb","brightWhite":"#eceff4"}'
            ;;
        "tokyo_night")
            # Tokyo Night theme
            echo '{"background":"#1a1b26","foreground":"#c0caf5","cursor":"#c0caf5","cursorAccent":"#1a1b26","selectionBackground":"#33467c","selectionForeground":"#c0caf5","black":"#15161e","red":"#f7768e","green":"#9ece6a","yellow":"#e0af68","blue":"#7aa2f7","magenta":"#bb9af7","cyan":"#7dcfff","white":"#a9b1d6","brightBlack":"#414868","brightRed":"#f7768e","brightGreen":"#9ece6a","brightYellow":"#e0af68","brightBlue":"#7aa2f7","brightMagenta":"#bb9af7","brightCyan":"#7dcfff","brightWhite":"#c0caf5"}'
            ;;
        "one_dark")
            # One Dark theme
            echo '{"background":"#282c34","foreground":"#abb2bf","cursor":"#528bff","cursorAccent":"#282c34","selectionBackground":"#3e4451","selectionForeground":"#abb2bf","black":"#282c34","red":"#e06c75","green":"#98c379","yellow":"#e5c07b","blue":"#61afef","magenta":"#c678dd","cyan":"#56b6c2","white":"#abb2bf","brightBlack":"#5c6370","brightRed":"#e06c75","brightGreen":"#98c379","brightYellow":"#e5c07b","brightBlue":"#61afef","brightMagenta":"#c678dd","brightCyan":"#56b6c2","brightWhite":"#ffffff"}'
            ;;
        "solarized_dark")
            # Solarized Dark theme
            echo '{"background":"#002b36","foreground":"#839496","cursor":"#839496","cursorAccent":"#002b36","selectionBackground":"#073642","selectionForeground":"#93a1a1","black":"#073642","red":"#dc322f","green":"#859900","yellow":"#b58900","blue":"#268bd2","magenta":"#d33682","cyan":"#2aa198","white":"#eee8d5","brightBlack":"#002b36","brightRed":"#cb4b16","brightGreen":"#586e75","brightYellow":"#657b83","brightBlue":"#839496","brightMagenta":"#6c71c4","brightCyan":"#93a1a1","brightWhite":"#fdf6e3"}'
            ;;
        "solarized_light")
            # Solarized Light theme
            echo '{"background":"#fdf6e3","foreground":"#657b83","cursor":"#657b83","cursorAccent":"#fdf6e3","selectionBackground":"#eee8d5","selectionForeground":"#586e75","black":"#073642","red":"#dc322f","green":"#859900","yellow":"#b58900","blue":"#268bd2","magenta":"#d33682","cyan":"#2aa198","white":"#eee8d5","brightBlack":"#002b36","brightRed":"#cb4b16","brightGreen":"#586e75","brightYellow":"#657b83","brightBlue":"#839496","brightMagenta":"#6c71c4","brightCyan":"#93a1a1","brightWhite":"#fdf6e3"}'
            ;;
        "gruvbox_dark")
            # Gruvbox Dark theme
            echo '{"background":"#282828","foreground":"#ebdbb2","cursor":"#ebdbb2","cursorAccent":"#282828","selectionBackground":"#504945","selectionForeground":"#ebdbb2","black":"#282828","red":"#cc241d","green":"#98971a","yellow":"#d79921","blue":"#458588","magenta":"#b16286","cyan":"#689d6a","white":"#a89984","brightBlack":"#928374","brightRed":"#fb4934","brightGreen":"#b8bb26","brightYellow":"#fabd2f","brightBlue":"#83a598","brightMagenta":"#d3869b","brightCyan":"#8ec07c","brightWhite":"#ebdbb2"}'
            ;;
        *)
            # Default to Breeze
            echo '{"background":"#232627","foreground":"#fcfcfc","cursor":"#fcfcfc","cursorAccent":"#232627","selectionBackground":"#3daee9","selectionForeground":"#fcfcfc","black":"#232627","red":"#ed1515","green":"#11d116","yellow":"#f67400","blue":"#1d99f3","magenta":"#9b59b6","cyan":"#1abc9c","white":"#fcfcfc","brightBlack":"#7f8c8d","brightRed":"#c0392b","brightGreen":"#1cdc9a","brightYellow":"#fdbc4b","brightBlue":"#3daee9","brightMagenta":"#8e44ad","brightCyan":"#16a085","brightWhite":"#ffffff"}'
            ;;
    esac
}

THEME_JSON=$(get_theme "${TERMINAL_THEME}")

# ==============================================================================
# Set Node.js memory limit (hardcoded for stability)
# ==============================================================================

export NODE_OPTIONS="--max-old-space-size=512"

bashio::log.info "Starting ttyd on port 8099..."

# ==============================================================================
# Start ttyd with configured options
# ==============================================================================

exec ttyd \
    -W \
    -p 8099 \
    -t fontSize="${FONT_SIZE}" \
    -t fontFamily="'Hack', 'DejaVu Sans Mono', 'Liberation Mono', 'Noto Sans Mono', monospace" \
    -t fontWeight=normal \
    -t fontWeightBold=bold \
    -t letterSpacing=0 \
    -t lineHeight=1.0 \
    -t cursorBlink="${CURSOR_BLINK}" \
    -t cursorStyle="${CURSOR_STYLE}" \
    -t drawBoldTextInBrightColors=true \
    -t minimumContrastRatio=4.5 \
    -t scrollback=10000 \
    -t "theme=${THEME_JSON}" \
    /usr/local/bin/opencode-session.sh
